# debian/postinst
#!/bin/sh
set -e

GROUP="tcall"
BINDIR="/usr/local/bin"
RUNDIR="/var/run/tcall"
LOGDIR="/var/log/tcall"
ETCDIR="/etc/tcall"

# group
getent group "$GROUP" >/dev/null || groupadd "$GROUP"

# add installing admin
if [ -n "$SUDO_USER" ]; then
  usermod -aG "$GROUP" "$SUDO_USER" || true
fi

# runtime dirs
mkdir -p "$RUNDIR" "$LOGDIR" "$ETCDIR"
chown root:"$GROUP" "$RUNDIR" "$LOGDIR"
chmod 770 "$RUNDIR" "$LOGDIR"

touch "$RUNDIR/user_request" "$RUNDIR/pending_approval" "$LOGDIR/audit.log"
chown root:"$GROUP" "$RUNDIR/"* "$LOGDIR/audit.log"
chmod 660 "$RUNDIR/"* "$LOGDIR/audit.log"

# policy
cp /etc/tcall/policy.txt /etc/tcall/policy.txt 2>/dev/null || true

# binaries
chmod 755 "$BINDIR"/tcall*
chown root:root "$BINDIR"/tcall*

# profile enforcement
chmod 755 /etc/profile.d/tcall.sh

# systemd
systemctl daemon-reload || true
systemctl enable tcall.service || true
systemctl restart tcall.service || true

exit 0
