# debian/prerm
#!/bin/sh
set -e

SERVICE="tcall.service"

systemctl stop "$SERVICE" || true
systemctl disable "$SERVICE" || true

exit 0
